# supervisor.conf
[program:ohrkanal-podcast]
command=docker-compose -f production.yml up
directory=/home/ohrkanal/ohrkanal
redirect_stderr=true
autostart=true
autorestart=true
priority=10
